# scripts/gt_energy.cron
#| MINUTE |	HOUR | DAY OF MONTH | MONTH | DAY OF WEEK | CMD#
   *5		*	   *			  *		  0-4				/Users/Shared/giant_theremin/scripts/mktimes
